# roles/infra_cloudflared/templates/cloudflared.service.j2
[Unit]
Description=Cloudflare Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ cloudflared_user }}
Group={{ cloudflared_group }}
ExecStart=/usr/local/bin/cloudflared tunnel --config {{ cloudflared_config_dir }}/config.yml run
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ cloudflared_config_dir }}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
